# server/config/rc.local
#!/bin/sh -e

# after dhcp
systemctl start avahi-alias@mqtt.local

exit 0
