# server/config/rc.local
#!/bin/sh -e

# after dhcp
if [ -f /opt/battery_pack/server/scripts/start ]; then
        /bin/sh /opt/battery_pack/server/scripts/start &
fi

exit 0
